# ansible/roles/sliceit_app/templates/supervisor_template.conf.j2
[program:{{ app_name }}]
command={{ app_dest }}/{{ app_dir }}/start.sh
autostart=true
autorestart=true
